# supervisord.conf
[supervisord]
nodaemon=true

[program:log.io]
command=log.io-server
autostart=true
autorestart=true